# backend/texchat.supervisord.conf


[program:texchat]
command=NODE_ENV=production node backend/server.js
directory=/var/www/texchat

startretries=10
stdout_logfile_maxbytes=5MB
priority=10
redirect_stderr=true
